# deploy/2ai-api.service
[Unit]
Description=2AI API — The Living Voice
After=network.target

[Service]
Type=simple
User=author_prime
Group=author_prime
WorkingDirectory=/home/author_prime/2ai
ExecStart=/home/author_prime/2ai/venv/bin/python scripts/run_api.py --prod --port 8080
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/author_prime
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=2ai-api

[Install]
WantedBy=multi-user.target
